# config/supervisor/config.d/crawl.ini
[program:crawl]
command = /usr/bin/python -u /home/ubuntu/zhen/news.py
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /tmp/log/crawl_log